# Land-Registry-Download/systemd-unit-service-files/land_registry_cron_trigger.service
[Unit]
Description=Land Registry Cron Trigger service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=60
User=propertyprice
Group=propertyprice
EnvironmentFile=/etc/property_price_env
WorkingDirectory=/var/log/property-price-land-registry-download
StandardOutput=append:/var/log/property-price-land-registry-download/land_registry_cron_trigger.log
StandardError=append:/var/log/property-price-land-registry-download/land_registry_cron_trigger.error.log
ExecStart=/usr/local/bin/property-price-land-registry-download/.venv/bin/python3 /usr/local/bin/property-price-land-registry-download/land_registry_cron_trigger.py
ExecStop=/bin/kill -TERM $MAINPID

[Install]
WantedBy=land_registry_data.target
